# 8765_Node.js_Listings/listings/Kapitel25/Listing25_3/express-app.service
[Unit]
Description=Express Application
After=network.target

[Service]
User=node
WorkingDirectory=/srv/node
ExecStart=npm start

[Install]
WantedBy=multi-user.target
